# etc/artyfarty.timer
[Timer]
OnCalendar=*:0/30:0
Persistent=false
AccuracySec=1s
