# etc/artyfarty.timer
[Timer]
OnCalendar=00/2:0:0
Persistent=false
AccuracySec=1s
